# 3ºCurso/Primer Cuatri/ISE/Practicas/P3/mon-raid/mon_raid.service
[Unit]
Description=Monitor RAID status

[Service]
Type=simple
ExecStart=/home/gabrieljp/mon-raid.py
